# backend/docker/php/crontab
* * * * * cd /opt &&  php ./bin/console documents:delete --env=dev
# Mandatory blank line